# labs/sleep.service
[Unit]
After=network-online.target


[Service]
ExecStartPre=/usr/bin/sleeplog.sh
ExecStart=/usr/bin/rsleep.sh
Restart=always
